# bin/tempserver-remote.service
[Unit]
Description=Tempremote deamon
[Service]
User=alex
ExecStart=java -jar release.jar
WorkingDirectory=/home/alex/remote_java_beta1
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target